# arkteos/root/etc/services.d/arkteos/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running arkteos addon"

# Declare variables
declare password_mqtt
declare log_level

# Get all arguments
password_mqtt=$(bashio::config 'password_mqtt')
export PASSWORDMQTT="${password_mqtt}"

log_level=$(bashio::config 'log_level')
export LOGLEVEL="${log_level}"

exec python3 /arkteos/arkteos.py
