# teleinfo/root/etc/services.d/teleinfo/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running teleinfo addon"

# Declare variables
declare port_usb

# Get all arguments
port_usb=$(bashio::config 'port_usb')
export USBPORT="${port_usb}"

exec python3 /teleinfo/teleinfo_standard.py
