# teleinfo/root/etc/services.d/arkteos/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Running Arkteos addon"

exec python3 -u /arkteos/arkteos.py
